# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:streamer]
command = python stream.py
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0